# .config/systemd/user/status-notifier-ready.service
[Unit]
Description=StatusNotifier tray readiness
After=waybar.service
PartOf=waybar.service

[Service]
Type=oneshot
ExecStart=/usr/bin/gdbus wait --session --timeout 5 org.kde.StatusNotifierWatcher
RemainAfterExit=yes
