# .config/systemd/user/app-@autostart.service.d/50-status-notifier.conf
[Unit]
Wants=status-notifier-ready.service
After=status-notifier-ready.service
